# agents/crons/session-sync/com.monorepo-starter.session-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.monorepo-starter.session-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/agents/crons/session-sync/script.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/memory/crons/session-sync.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/memory/crons/session-sync.stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
